# etc/cron.d/backup
0 14 * * 2,5,7	root	/opt/runlevel_backup/backup.sh > /dev/null
